# endpoint/agent-endpoint.service
[Unit]
Description=Koine agent answer-endpoint (full-context daemon: human/ops channels)
Documentation=https://github.com/dwilkin/koine (endpoint/README.md)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# --- non-secret config; the bearer token is read from your secret store by run.sh, never here ---
# REFERENCE UNIT: replace the <placeholders> per deployment (a domain keeps its filled-in copy
# in its own ops repo — e.g. the wilkin lab overlay lives in its private ops repo, not here).
Environment=AGENT_NAME=<agent>
Environment=WORKDIR=%h/<project-dir-with-CLAUDE.md>
Environment=CLAUDE_BIN=%h/.local/bin/claude
# run.sh reads the bearer(s) from your secret store — pin its address here if the unit's
# environment doesn't otherwise carry it (e.g. Environment=VAULT_ADDR=https://<vault-host>:8200)
Environment=ENDPOINT_BIND=0.0.0.0:8090
Environment=MODEL=sonnet
# Human-channel answers may get a smarter model (chat should feel like a real conversation);
# peer A2A stays on the cheaper MODEL.
Environment=MODEL_HUMAN=opus
Environment=ANSWER_TIMEOUT=180
Environment=ANSWER_TIMEOUT_HUMAN=300
# Chat-execution parity: make the human channel a real control channel — live-session
# permission mode; guard hooks stay the hard floor; confirm-gated work still waits for the
# human's in-chat approval via the pending-actions ledger. Remove to keep chat advisory-only.
Environment=PERMISSION_MODE_HUMAN=bypassPermissions
Environment=MAX_CONCURRENCY=2
# Peer-path security alerts (redaction/tripwire hits) go to the human via your domain's
# notify helper (any executable taking the alert text as argv[1]). Empty/unset = audit-only.
Environment=ALERT_CMD=<path-to-your-notify-helper>
ExecStart=%h/koine/endpoint/run.sh
Restart=on-failure
RestartSec=5
# hardening (best-effort for a --user service)
NoNewPrivileges=true

[Install]
WantedBy=default.target
